# package/write-ram-logs/write_ram_logs.service
[Unit]
Description=Write tmpfs to disk

[Service]
ExecStart=/path/to/write_ram_logs.sh
