# package/write-ram-logs/write-ram-logs.service
[Unit]
Description=Write tmpfs to disk

[Service]
ExecStart=/usr/bin/write-ram-logs.sh

[Install]
WantedBy=multi-user.target
